# private_dot_config/exact_systemd/exact_user/shairport-sync.service
[Unit]
Description=Shairport Sync - AirPlay Audio Receiver

[Service]
ExecStart=/usr/bin/shairport-sync
Restart=always

[Install]
WantedBy=default.target
