# private_dot_config/exact_systemd/exact_user/wallpaper.service.tmpl
[Unit]
Description=Change sway background

[Service]
Type=oneshot
ExecStart={{ .chezmoi.homeDir }}/.local/bin/wallpaper
