# private_dot_config/exact_systemd/exact_user/waynergy.service
[Unit]
Description=waynergy

[Service]
ExecStart=/usr/bin/waynergy
Restart=on-failure

[Install]
WantedBy=default.target
