# private_dot_config/exact_autostart/1password.desktop
[Desktop Entry]
Name=1Password
Exec=1password --silent
Type=Application
StartupNotify=false
