# private_dot_config/exact_autostart/autotiling.desktop
[Desktop Entry]
Name=Autotiling
Exec=autotiling
Type=Application
OnlyShowIn=sway;
StartupNotify=false
