# private_dot_config/exact_autostart/clipman-restore.desktop
[Desktop Entry]
Name=Mako
Exec=clipman restore --notify
Type=Application
OnlyShowIn=sway;
StartupNotify=false
