# private_dot_config/exact_autostart/clipman.desktop
[Desktop Entry]
Name=Mako
Exec=wl-paste -t text --watch clipman store --no-persist --notify
Type=Application
OnlyShowIn=sway;
StartupNotify=false
