# private_dot_config/exact_autostart/ibus-daemon.desktop
[Desktop Entry]
Name=IBus Daemon
Exec=ibus-daemon -drx --panel /usr/lib/ibus/ibus-ui-gtk3
Type=Application
OnlyShowIn=sway;
StartupNotify=false
