# private_dot_config/exact_autostart/mako.desktop
[Desktop Entry]
Name=Mako
Exec=mako
Type=Application
OnlyShowIn=sway;
StartupNotify=false
