# private_dot_config/exact_autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Exec=udiskie --smart-tray --appindicator
Type=Application
OnlyShowIn=sway;
StartupNotify=false
