# private_dot_config/exact_autostart/volume-bar.desktop
[Desktop Entry]
Name=Volume Bar
Exec=/home/jrubin/bin/volume-bar
Type=Application
OnlyShowIn=sway;
StartupNotify=false
