# private_dot_config/exact_systemd/exact_user/pulseaudio.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio --daemonize=no --high-priority --realtime --exit-idle-time=-1
