# private_dot_config/systemd/user/wallpaper.timer
[Unit]
Description=Change sway background every hour

[Timer]
OnCalendar=*:0

[Install]
WantedBy=timers.target
