# private_dot_config/systemd/user/wallpaper.service
[Unit]
Description=Change sway background

[Service]
Type=oneshot
ExecStart=/home/jrubin/bin/wallpaper
